# deploy/admarketplace-backend.service
[Unit]
Description=AdMarketplace API Backend (uvicorn)
After=network.target postgresql.service
Wants=postgresql.service

[Service]
Type=simple
User=admarket
Group=admarket
WorkingDirectory=/opt/admarketplace/backend_py

ExecStart=/opt/admarketplace/backend_py/.venv/bin/uvicorn app.main:app --host 127.0.0.1 --port 3100
Restart=always
RestartSec=10

Environment=PYTHONUNBUFFERED=1
EnvironmentFile=/opt/admarketplace/backend_py/.env

[Install]
WantedBy=multi-user.target
